# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:backend]
command=/bin/sh -c "cd /app/backend && python init_db.py && exec uvicorn main:app --host 0.0.0.0 --port 8000"
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:frontend]
command=/usr/sbin/nginx -g "daemon off;"
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
